# chapter-5/myblog_web.conf
[program:myblog_app]
command=/usr/bin/python /srv/mezzanine/manage.py runserver
autostart=true
autorestart=unexpected
stopwaitsecs=10
stopasgroup=true
killasgroup=true
user=mezzanine